# crontab
# blacklisting from 0:00-5:00

30 23 * * * /home/pi/safeEnable
0 5 * * * /home/pi/safeDisable

# update gravity daily

0 0 * * * pihole -g

# update pihole every sunday at 0:00

0 0 * * 0 pihole -up
