# crontab
# example crontab use
# blacklist in effect 0:00-5:00

0 0 * * * /home/pi/safeEnable
0 5 * * * /home/pi/safeDisable
